# infrastructure/systemd/theatre-backup.timer
[Unit]
Description=Theatre Backup to GCS Timer
Documentation=https://github.com/szmyty/theatre

[Timer]
# Run daily at 3:00 AM UTC
OnCalendar=*-*-* 03:00:00
# Randomize start time by up to 30 minutes to avoid thundering herd
RandomizedDelaySec=1800
# Run immediately if missed while system was offline
Persistent=true

[Install]
WantedBy=timers.target
